# config/chroot_local-includes/lib/systemd/system/tor@default.service.d/writable-etc-tor.conf
[Service]

# SAVECONF requires tor to be allowed to modify its own configuration
ProtectSystem=yes
ReadWriteDirectories=-/etc/tor
