# config/chroot_local-includes/lib/systemd/system/console-setup.service.d/override.conf
[Unit]
# Workaround tails#18636 (https://bugs.debian.org/846256)
After=systemd-tmpfiles-setup.service
