# config/chroot_local-includes/lib/systemd/system/tails-persistent-storage.service
[Service]
Type=dbus
BusName=org.boum.tails.PersistentStorage
ExecStartPre=/usr/local/lib/persistent-storage/pre-start
ExecStart=/usr/local/lib/tps
