# config/chroot_local-includes/lib/systemd/system/tails-additional-software-upgrade.path
[Unit]
Description=Trigger upgrade of Additional Software Packages
Documentation=https://tails.net/contribute/design/persistence/
After=tails-additional-software-install.service
After=tails-wait-until-tor-has-bootstrapped.service
ConditionFileNotEmpty=/live/persistence/TailsData_unlocked/live-additional-software.conf
ConditionPathIsMountPoint=/var/lib/apt/lists
ConditionPathIsMountPoint=/var/cache/apt/archives

[Path]
PathExists=/run/live-additional-software/installed
