# config/chroot_local-includes/usr/lib/systemd/user/org.gnome.Shell@wayland.service.d/override.conf
[Service]
ExecStartPre=/usr/local/lib/gnome-shell-dump-environment
